# systemd/wifi-fix.service
[Unit]
Description=WiFi Retry on Boot
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c "sleep 10; nmcli con up SkyNet2G 2>/dev/null || true"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
